# config/com.claude.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.sync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/eric/ccproj/claude_sync_startup.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/eric/.claude-sync.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/eric/.claude-sync.err</string>
    
    <key>StartInterval</key>
    <integer>3600</integer>
    
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>